# examples/init.d/jhttpd.in
#! /bin/bash
#
# jhttpd          Start/Stop a http daemon.
#
# chkconfig: 345 90 60
# description: jhttpd is a simple http/1.0 server
# processname: ex_httpd
# pidfile: /var/run/jhttpd.pid
 
# Source function library.
. /etc/init.d/functions

# Needed to resolve sbindir ... *sigh*
prefix="@prefix@"
exec_prefix="@exec_prefix@"
localstatedir="@localstatedir@"

prog_name=jhttpd
prog_exe="@sbindir@/${prog_name}"

cntl_exe="@sbindir@/jcntl"

MIME_EXTRAS_FILE="@libdir@/vstr-@VERSION@/examples/data/mime_types_extra.txt"

FILE_BASE="${localstatedir}/run/${prog_name}"
JHTTPD_NUM=1
JHTTPD1_DIR=${localstatedir}/www/html
[ -r /etc/sysconfig/jhttpd ] && . /etc/sysconfig/jhttpd

GLOB_I_ARGS="--daemon --mime-types-xtra $MIME_EXTRAS_FILE"

function jhttpd_dirs()
{
          CNTL_FILE="${FILE_BASE}$i.cntl"
          PID_FILE="${FILE_BASE}$i.pid"
          U__CALL_CHROOT="JHTTPD${i}_CHROOT"
          U__CALL_DIR="JHTTPD${i}_DIR"
	  if [ "x${!U__CALL_CHROOT}" = "x" ]; then
	    REAL_DIR="${!U__CALL_DIR:-/}"
	  else
	    REAL_DIR="${!U__CALL_CHROOT}${!U__CALL_DIR:+/}${!U__CALL_DIR}"
	  fi
}

function start()
{
        RETVAL=0
        i=1
	while [ $RETVAL -eq 0 -a $i -le $JHTTPD_NUM ];
        do
          jhttpd_dirs
          U__CALL_ARGS="JHTTPD${i}_ARGS"
          I_ARGS="$GLOB_I_ARGS --daemon --cntl-file ${CNTL_FILE} --pid-file ${PID_FILE}"
          U_ARGS="${!U__CALL_ARGS} ${!U__CALL_CHROOT:+--chroot} ${!U__CALL_CHROOT} ${!U__CALL_DIR:-/}"
          echo -n $"Starting $prog_name(${REAL_DIR}): "
          daemon --check $prog_name$i $prog_exe $I_ARGS $U_ARGS
          RETVAL=$?
          echo
          i=$(($i + 1))
        done
        [ $RETVAL -eq 0 ] && touch /var/lock/subsys/$prog_name
        return $RETVAL
}

function stop()
{
        RETVAL=0
        i=1
	while [ $RETVAL -eq 0 -a $i -le $JHTTPD_NUM ];
        do
          jhttpd_dirs
          echo -n $"Shutting down $prog_name(${REAL_DIR}): "
          killproc $prog_name$i
          RETVAL=$?
	  echo
          i=$(($i + 1))
        done
        [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/$prog_name
        return $RETVAL
}

function softstop()
{
        RETVAL=0
        i=1
	while [ $RETVAL -eq 0 -a $i -le $JHTTPD_NUM ];
        do
          jhttpd_dirs
          echo -n $"Soft shutdown $prog_name(${REAL_DIR}): "
          $cntl_exe $CNTL_FILE -e close > /dev/null
          RETVAL=$?
          [ $RETVAL -eq 0 ] && success $"$prog_name$i soft shutdown" || \
	     failure $"$prog_name$i soft shutdown"
          [ $RETVAL -eq 0 ] && rm -- "${PID_FILE}"
	  echo
          i=$(($i + 1))
        done
        [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/$prog_name
        return $RETVAL
}

function restart()
{
        stop; start;
}

function softrestart()
{
        softstop; start;
}

function mystatus()
{
        # First try "pidof"
        pids=`pidof -o $$ -o $PPID -o %PPID -x $1`
        if [ -n "$pids" ]; then
                if [ -f /var/run/$1$2.pid ] ; then
                        read pid < /var/run/$1$2.pid
			for j in $pids; do
			    if [ "$pid" = "$j" ]; then
                                echo $"$1$2 (pid $pid) is running..."
                                return 0
			    fi
			done
                fi
        fi
 
        # Next try "/var/run/*.pid" files
        if [ -f /var/run/$1$2.pid ] ; then
                read pid < /var/run/$1$2.pid
                if [ -n "$pid" ]; then
                        echo $"$1$2 dead but pid file exists [$pid]"
			return 1
                fi
        fi
        # See if /var/lock/subsys/${base} exists
        if [ -f /var/lock/subsys/$1 ]; then
                echo $"$1$2 dead but subsys locked"
                return 2
        fi
        echo $"$1 is stopped"
        return 3
}

function statuses()
{
        i=1
	while [ $i -le $JHTTPD_NUM ];
        do
          jhttpd_dirs
          mystatus $prog_name $i
          RETVAL=$?
          [ $RETVAL -eq 0 ] && $cntl_exe $CNTL_FILE -e status
          i=$(($i + 1))
        done
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart)
        restart
        ;;
  softstop)
        softstop
        ;;
  softrestart)
        softrestart
        ;;
  graceful)
        softstop
        ;;
  reload)
        softrestart
        ;;
  status)
        statuses
        ;;
  condrestart)
        [ -f /var/lock/subsys/${prog_name}1 ] && softrestart || :
        ;;
  *)
	echo $"Usage: $0 {start|stop|restart|softstop|softrestart|reload|status|condrestart}"
        exit 1
esac
                                                                             
exit $?
