# examples/init.d/jhttpd.in
#! /bin/bash
#
# jhttpd          Start/Stop a http daemon.
#
# chkconfig: 345 90 60
# description: jhttpd is a simple http/1.0 server
# processname: ex_httpd
# pidfile: /var/run/jhttpd.pid
 
# Source function library.
. /etc/init.d/functions

# Needed to resolve sbindir ... *sigh*
prefix="@prefix@"
exec_prefix="@exec_prefix@"

prog_name=jhttpd
prog_exe="@sbindir@/${prog_name}"

I_ARGS="--daemon --pid-file /var/run/${prog_name}.pid"
JHTTPD_DIR=/var/www/html
[ -r /etc/sysconfig/jhttpd ] && . /etc/sysconfig/jhttpd

function start()
{
        echo -n $"Starting $prog_name: "
        daemon --check $prog_name $prog_exe $I_ARGS $JHTTPD_ARGS $JHTTPD_DIR
        RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && touch /var/lock/subsys/$prog_name
        return $RETVAL
}

function stop()
{
        echo -n $"Shutting down $prog_name: "
        killproc $prog_name
        RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/$prog_name
        return $RETVAL
}

function restart()
{
        stop; start;
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart)
        restart
        ;;
  reload)
        restart
        ;;
  status)
        status $prog_name
        ;;
  condrestart)
        [ -f /var/lock/subsys/$prog_name ] && restart || :
        ;;
  *)
	echo $"Usage: $0 {start|stop|restart|condrestart|status}"
        exit 1
esac
                                                                             
exit $?
